# infrastructure/systemd/pincer.service
[Unit]
Description=Pincer Agent Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=micelio
Group=micelio
WorkingDirectory=/home/micelio/git/Pincer
Environment=MIX_ENV=prod
EnvironmentFile=-/etc/pincer/pincer.env
ExecStart=/usr/bin/env mix pincer.server telegram
Restart=always
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/micelio/git/Pincer /home/micelio/.mix /home/micelio/.cache/rebar3
UMask=0077
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pincer

[Install]
WantedBy=multi-user.target
